# ansible/files/dnsmasq/dnsmasq_adblock.timer
[Unit]
Description = timer adblock

[Timer]
OnUnitActiveSec=1d

[Install]
WantedBy=dnsmasq_adblock.target
