# ansible/files/gsnet_scanner/gsnet_scanner.service
[Unit]
Description = gsnet scanner

[Service]
ExecStart = /usr/local/bin/gsnet_scanner
Type = simple

[Install]
WantedBy=multi-user.target
